# overlay/etc/crontabs/nginx
*/15 * * * * /usr/bin/php /var/www/app/update.php --feeds
