# overlay/etc/services.d/update-feeds/run
#!/bin/execlineb -P
s6-setuidgid nginx php7 /var/www/app/update_daemon2.php > /dev/null 2>&1
